# systemd/sabori-flow.service.template
[Unit]
Description=sabori-flow worker

[Service]
Type=oneshot
ExecStart=__EXEC_START__
Environment=PATH=__PATH__
StandardOutput=append:__LOG_DIR__/systemd_stdout.log
StandardError=append:__LOG_DIR__/systemd_stderr.log
